# debian/incus.prerm
#!/bin/sh
set -e

case "$1" in
    remove)
        systemctl disable --now incus-startup.service
        systemctl disable --now incus.service incus.socket
        systemctl disable --now incus-lxcfs.service
    ;;

    upgrade|deconfigure|failed-upgrade)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

#DEBHELPER#

exit 0
